# crontab
## imperial-notifier
*/5 * * * * cd /web/imperial-notifier && ENV=prod /home/web/.local/bin/poetry run python3.12 -m check >> log/notifier.log 2>&1
